# root/etc/s6-overlay/s6-rc.d/init-mod-swag-auto-proxy-setup/finish
#!/usr/bin/with-contenv bash

sed -i '/\/app\/auto-proxy-external.sh/d' /etc/crontabs/root
rm -rf /config/nginx/proxy-confs/auto-proxy*.conf